# deploy/voice-server.service
[Unit]
Description=Voice Server (Node.js Express + WS + STUN)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__VOICE_SSH_USER__
Group=__VOICE_SSH_USER__
WorkingDirectory=/opt/voice-server
EnvironmentFile=/etc/voice-server.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
# Reasonable hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
