# deploy/voice-edge-cert-renew.service
[Unit]
Description=Voice-server edge self-signed TLS cert renewal check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/voice-edge/scripts/cert-renew.sh
StandardOutput=journal
StandardError=journal
